# blocklists_ipset.service
[Unit]
Description=Update the blocklist from blocklists.de

[Service]
ExecStart=/root/blocklists-ipset.py